# api/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:api]
command=/usr/local/bin/api

[program:kratos]
command=/usr/local/bin/kratos
