# api/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:api]
command=/usr/local/bin/api

[program:kratos]
command=/usr/local/bin/kratos serve
